# init.d/dbus
#!/sbin/openrc-run
name="System Message Bus"
description="Service for the D-Bus Daemon (system dbus)"

extra_started_commands="reload"

command="/usr/local/bin/dbus-daemon"
command_args="--system --nofork --nopidfile --syslog-only ${command_args:-}"
command_background="yes"
pidfile="/var/run/$RC_SVCNAME.pid"

depend(){
	need localmount
	need FILESYSTEMS
}

start_pre(){
	checkpath -d -m755 -o root:messagebus /var/run/dbus || return 1

	/usr/local/bin/dbus-uuidgen --ensure="${machine_id:-/etc/machine-id}"
}

stop_post(){
	[ ! -S /var/run/dbus/system_bus_socket ] || rm -f /var/run/dbus/system_bus_socket
}

reload(){
	ebegin "Reloading $name configuration"
	/usr/local/bin/dbus-send --print-reply --system --type=method_call --dest=org.freedesktop.DBus org.freedesktop.DBus.ReloadConfig > /dev/null
	eend $?
}
