# init.d/iovctl
#!/sbin/openrc-run
#
#

# PROVIDE: iovctl
# REQUIRE: FILESYSTEMS sysctl

. /etc/rc.subr
. /etc/netrc.conf

name="iovctl"
command="/usr/sbin/iovctl"

depend(){
	need FILESYSTEMS
	need sysctl
}

run_iovctl()
{
	local _f flag

	flag=$1
	for _f in ${iovctl_files} ; do
		if [ -r ${_f} ]; then
			${command} ${flag} -f ${_f} > /dev/null
		fi
	done
}

start()
{
	run_iovctl -C
}

stop()
{
	run_iovctl -D
}

