# init.d/musicpd
#!/sbin/openrc-run

name=$RC_SVCNAME
cfgfile="/etc/${RC_SVCNAME}/${RC_SVCNAME}.conf"
#supervisor=supervise-daemon
#USER=kavex
MPDCONF=/home/kavex/.config/mpd/mpd.conf

#command_user="${USER}"
command_args="--no-daemon ${MPDCONF}"
command="/usr/local/bin/musicpd"
pidfile="/var/run/${RC_SVCNAME}/${RC_SVCNAME}.pid"
command_background="yes"

depend(){
	need network
}

#start_pre(){
#	checkpath -d ${MPDCONF} 
#	${command_user} --mode 0775 /run/${RC_SVCNAME}
#	/var/mpd/log
#}
