# init.d/zfs
#!/sbin/openrc-run
#
#

# PROVIDE: zfs
# REQUIRE: zfsbe
# BEFORE: FILESYSTEMS var

. /etc/rc.subr
. /etc/netrc.conf

PATH=/sbin:/bin:/usr/sbin:/usr/bin:/usr/local/bin:/usr/local/sbin
export PATH

name="zfs"
desc="Mount and share ZFS datasets"
rcvar="zfs_enable"
#required_modules="zfs"

depend(){
	need zfsbe
	use FILESYSTEMS var
}

zfs_start_jail()
{
	if [ `$SYSCTL_N security.jail.mount_allowed` -eq 1 ]; then
		zfs mount -a
	fi
}

zfs_start_main()
{
	zfs mount -va
	zfs share -a
	if [ ! -r /etc/zfs/exports ]; then
		touch /etc/zfs/exports
	fi
}

start()
{
	if [ `$SYSCTL_N security.jail.jailed` -eq 1 ]; then
		zfs_start_jail
	else
		zfs_start_main
	fi
}

poststart()
{
	# Some of the keys to decrypt datasets are potentially stored on ZFS
	# datasets that just got mounted. Let's try to load those keys and
	# mount the datasets.
	if checkyesno zfskeys_enable; then
		/etc/rc.d/zfskeys start
		zfs_start
	fi
}

zfs_stop_jail()
{
	if [ `$SYSCTL_N security.jail.mount_allowed` -eq 1 ]; then
		zfs unmount -a
	fi
}

zfs_stop_main()
{
	zfs unshare -a
	zfs unmount -a
}

stop()
{
	if [ `$SYSCTL_N security.jail.jailed` -eq 1 ]; then
		zfs_stop_jail
	else
		zfs_stop_main
	fi
}

