# init.d/zpoolreguid
#!/sbin/openrc-run

# PROVIDE: zpoolreguid
# REQUIRE: zpool
# BEFORE: mountcritlocal
# KEYWORD: firstboot nojail

name="zpoolreguid"
description="Generate a new zpool GUID"

depend(){
	need zpool
	use localmount
}

start()
{
	local pool

	for pool in ${zpool_reguid}; do
		zpool reguid $pool
	done
}
